# mosquitto.service
[Unit]
Description=Mosquitto MQTT broker

[Service]
Type=OneShot
ExecStart=mosquitto -d

[Install]
WantedBy=basic.target
